# monitor/monitor-detector.service
[Unit]
Description=Monitor Detector Service
After=network.target monitor-qqbot.service
Wants=network.target

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/monitor
Environment="PATH=/opt/monitor/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/opt/monitor/.env
ExecStart=/opt/monitor/venv/bin/python img_detect.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitor-detector

[Install]
WantedBy=multi-user.target